# releng/airootfs/etc/systemd/system/multi-user.target.wants/coffee-wallpapers-iso.service
[Unit]
Description=Fixes Wallpapers
After=network.target.target

[Service]
Type=oneshot
ExecStart=/opt/wallpaper-updater-iso

[Install]
WantedBy=multi-user.target
