# coffeelinux-archiso mod/releng/airootfs/opt/coffee-injector.service
[Unit]
Description=Updates CoffeeLinux
After=network.target.target

[Service]
Type=oneshot
ExecStart=/opt/coffee-injector

[Install]
WantedBy=multi-user.target
